# pr0bot.service
[Unit]
Description=pr0gramm Discord Bot
After=network.target

[Service]
Type=simple
User=dein_benutzername
WorkingDirectory=/pfad/zu/deinem/bot/Pro bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
